# day1/myweb.service
[Unit]
Description=Mój webserwer na dokerze
After=docker.service
Requires=docker.service

[Service]
ExecStart=/usr/bin/docker start -a webserver
ExecStop=/usr/bin/docker stop -t 2 webserver
Restart=always

[Install]
WantedBy=default.target
